# supervisor.conf.d/supervisor-blog.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /code/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
[program:django]
command= /usr/bin/python /code/manage.py migrate
